# kerchunkd.service
[Unit]
Description=kerchunkd GMRS Repeater Controller
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kerchunkd -c /etc/kerchunk/kerchunk.conf -f
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=kerchunk
Group=kerchunk
WorkingDirectory=/etc/kerchunk

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/kerchunk /etc/kerchunk /tmp
ProtectHome=true

# Allow access to USB HID and audio devices
SupplementaryGroups=audio plugdev

[Install]
WantedBy=multi-user.target
