# debian/kerchunk.postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create kerchunk system user and group
        if ! getent group kerchunk >/dev/null 2>&1; then
            addgroup --system kerchunk
        fi
        if ! getent passwd kerchunk >/dev/null 2>&1; then
            adduser --system --ingroup kerchunk --home /var/lib/kerchunk \
                    --no-create-home --gecos "kerchunkd daemon" kerchunk
        fi

        # Add kerchunk user to audio and plugdev groups for HID/sound access
        usermod -a -G audio kerchunk 2>/dev/null || true
        usermod -a -G plugdev kerchunk 2>/dev/null || true

        # Create runtime directories
        install -d -o kerchunk -g kerchunk -m 0755 /var/lib/kerchunk
        install -d -o kerchunk -g kerchunk -m 0755 /var/lib/kerchunk/voicemail
        install -d -o kerchunk -g kerchunk -m 0755 /var/lib/kerchunk/recordings
        install -d -o kerchunk -g kerchunk -m 0755 /var/lib/kerchunk/cdr
        install -d -o kerchunk -g kerchunk -m 0750 /var/log/kerchunk
        install -d -o kerchunk -g kerchunk -m 0755 /run/kerchunk

        # Create TTS cache directory
        install -d -o kerchunk -g kerchunk -m 0755 /usr/share/kerchunk/sounds/cache
        install -d -o kerchunk -g kerchunk -m 0755 /usr/share/kerchunk/sounds/cache/tts

        # Install default config if none exists
        if [ ! -f /etc/kerchunk/kerchunk.conf ]; then
            cp /etc/kerchunk/kerchunk.conf.example /etc/kerchunk/kerchunk.conf
            chown kerchunk:kerchunk /etc/kerchunk/kerchunk.conf
            chmod 0640 /etc/kerchunk/kerchunk.conf
        fi

        # Install default AI system prompt if none exists. The file is
        # owned by the kerchunk user so the operator (or the web UI)
        # can edit it in place; mod_ai re-reads it on config reload.
        if [ ! -f /etc/kerchunk/system_prompt.md ] && \
           [ -f /etc/kerchunk/system_prompt.md.example ]; then
            cp /etc/kerchunk/system_prompt.md.example /etc/kerchunk/system_prompt.md
            chown kerchunk:kerchunk /etc/kerchunk/system_prompt.md
            chmod 0644 /etc/kerchunk/system_prompt.md
        fi

        # Ensure config directory is writable for users.conf updates from web UI
        chown kerchunk:kerchunk /etc/kerchunk
        chmod 0755 /etc/kerchunk

        # Reload systemd to pick up the new service unit
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
        fi
        ;;
esac

#DEBHELPER#

# Belt-and-braces restart: if the service is running, restart it on
# the new binary, regardless of what deb-systemd-helper thinks about
# its enabled state. Operators who enabled with plain `systemctl
# enable kerchunkd` (not via deb-systemd-helper) otherwise see
# "kerchunkd.service is a disabled or a static unit not running, not
# starting it." on every upgrade and have to manually restart to pick
# up the new build.
if [ -d /run/systemd/system ] && \
   systemctl is-active --quiet kerchunkd.service 2>/dev/null; then
    systemctl try-restart kerchunkd.service >/dev/null 2>&1 || true
fi

exit 0
